# suesskartoffel/etc/init/batadv-vis-mz.conf
# Starts batadv-vis
#

description     "batadv-vis"
author          "Tobias Hachmer <tobias@hachmer.de>"

start on started alfred-mz
stop on stopped alfred-mz

respawn
script
. /etc/default/alfred-mz
exec /usr/sbin/batadv-vis -i ${BATMANIF} -u ${UNIX_SOCKET} -s
end script
