# zitronengras/etc/init/batadv-vis-wi.conf
# Starts batadv-vis
#

description     "batadv-vis"
author          "Tobias Hachmer <tobias@hachmer.de>"

start on started alfred-wi
stop on stopped alfred-wi

respawn
script
. /etc/default/alfred-wi
exec /usr/sbin/batadv-vis -i ${BATMANIF} -u ${UNIX_SOCKET} -s
end script
